# worldmonitor/rootfs/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:redis]
priority=1
; RDB snapshots to /data (the one directory HA Supervisor actually persists
; across restarts/updates) instead of --save "" --dir /tmp — without this,
; every add-on restart wiped the entire seeded dataset and forced a full
; ~45-70 minute reseed across all ~140 scripts from empty. requirepass still
; regenerates fresh each boot (see entrypoint.sh) since it's just an auth
; credential, not part of the dataset — reloading old data under a new
; password is fine, nothing else caches this one.
command=redis-server --bind 127.0.0.1 --port 6379 --requirepass %(ENV_REDIS_PASSWORD)s --maxmemory 256mb --maxmemory-policy allkeys-lru --save 60 100 300 10 --dir /data/redis
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:redis-rest]
priority=2
command=node /app/redis-rest/redis-rest-proxy.mjs
directory=/app/redis-rest
environment=SRH_TOKEN="%(ENV_REDIS_TOKEN)s",SRH_CONNECTION_STRING="redis://:%(ENV_REDIS_PASSWORD)s@127.0.0.1:6379",PORT="8079"
autostart=true
autorestart=true
startretries=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
priority=3
command=/usr/sbin/nginx -c /tmp/nginx.conf -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:worldmonitor-api]
priority=3
command=node /app/local-api-server.mjs
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:seed-loop]
priority=4
command=/app/seed-loop.sh
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:ais-relay]
priority=4
command=/app/ais-relay-wrapper.sh
directory=/app/seed
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
